# deployment/systemd/celery.service.template
[Unit]
Description=Celery Service for GigPulse
After=network.target

[Service]
User={{DEPLOY_USER}}
Group=www-data
WorkingDirectory={{APP_PATH}}
Environment="PATH={{VENV_PATH}}/Scripts"
ExecStart={{VENV_PATH}}/Scripts/celery -A config worker --loglevel=info

[Install]
WantedBy=multi-user.target
